# templates/datamount.mount.j2
[Unit]
Description=Media Storage at {{ data_mount_root }}
After=network-online.target firewalld.service
Wants=network-online.target

[Mount]
What={% for host in groups['mons'] %}{{ hostvars[host]['ansible_default_ipv4']['address'] }}:6789{% if not loop.last %},{% endif %}{% endfor %}:/
Where={{ data_mount_root }}
Type=ceph
Options=name={{  cephfs_user }},secret={{ vault_cephfs_user_key }},noatime,_netdev

[Install]
WantedBy=multi-user.target
