# templates/datamount.mount.j2
[Unit]
Description=Media Storage at {{ cephfs_mount.mount }}
After=network-online.target firewalld.service
Wants=network-online.target

[Mount]
What={% for host in cephfs_mount.mds %}{{ host }}{% if not loop.last %},{% endif %}{% endfor %}:/
Where={{ cephfs_mount.mount }}
Type=ceph
Options=name={{ cephfs_mount.user }},secret={{ cephfs_mount.secret }},noatime,_netdev

[Install]
WantedBy=multi-user.target
